# sys-cluster/spark/files/1.1.0/spark-master
#!/sbin/runscript

test -f /etc/spark/spark-env.sh && . /etc/spark/spark-env.sh

pidfile="/var/run/spark-master.pid"

check_config() {
	if [ ! -f /etc/spark/spark-env.sh ]
	then
		eerror "Please create /etc/spark/spark-env.sh"
		eerror "A template was installed at /etc/spark/spark-env.sh.template"
		return 1
	fi
	if [ -z $SPARK_MASTER_IP ]
	then
		eerror "You must set SPARK_MASTER_IP in /etc/spark/spark-env.sh"
		eerror "The default is \"\`hostname\`\""
		return 1
	fi
	if [ -z $SPARK_MASTER_PORT ]
	then
		eerror "You must set SPARK_MASTER_PORT"
		eerror "Default is 7077"
		return 1
	fi
	if [ -z $SPARK_MASTER_WEBUI_PORT ]
	then
		eerror "You must set SPARK_MASTER_WEBUI_PORT"
		eerror "Default is 8080"
		return 1
	fi
	return 0
}

start() {
	ebegin "Starting spark-master"
	check_config || return $?
	nohup /usr/lib/spark-1.1.0/bin/spark-class \
		org.apache.spark.deploy.master.Master \
		--ip $SPARK_MASTER_IP --port $SPARK_MASTER_PORT \
		--webui-port $SPARK_MASTER_WEBUI_PORT \
		> /dev/null 2>&1 < /dev/null &
	procpid=$!
	echo $procpid > $pidfile
	eend $? "Failed to start spark-master"
}

stop() {
	ebegin "Stopping spark-master"
	check_config || return $?

	local pid=`cat $pidfile`
	if [ -z `kill -0 $pid` ]
	then
		kill `cat $pidfile`
		sleep 5
		ps `cat $pidfile` > /dev/null
		if [ "$?" = "0" ]
		then
			echo "Didn't gracefully shutdown within 5 seconds, killing..."
			kill -9 `cat $pidfile`
		fi
		rm $pidfile
	else
		echo "Spark wasn't running"
	fi
	eend $? "Failed to stop spark-master"
}
